# deploy/sync-dashboard.service
# streamlit-dashboard-deploy FR-04 — 대시보드 데이터 동기화 (oneshot)
[Unit]
Description=auto-stock dashboard data sync (curate + push dashboard-data branch)

[Service]
Type=oneshot
User=auto-stock
WorkingDirectory=/opt/auto-stock
EnvironmentFile=/opt/auto-stock/.env
ExecStart=/opt/auto-stock/venv/bin/python scripts/sync_dashboard_data.py
