# deploy/sync-dashboard.service
# streamlit-dashboard-deploy FR-04 — 대시보드 데이터 동기화 (oneshot)
[Unit]
Description=SentiQuant dashboard data sync (curate + push dashboard-data branch)

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/auto_stock
EnvironmentFile=/home/ubuntu/auto_stock/.env
ExecStart=/home/ubuntu/auto_stock/venv/bin/python scripts/sync_dashboard_data.py
