# deploy/watchdog.service
# live-scheduler-deploy FR-09 — 워치독 (oneshot, timer가 주기 호출)
[Unit]
Description=auto-stock watchdog — heartbeat stale 감지 시 재시작/알림

[Service]
Type=oneshot
User=auto-stock
WorkingDirectory=/opt/auto-stock
EnvironmentFile=/opt/auto-stock/.env
ExecStart=/opt/auto-stock/venv/bin/python scripts/watchdog_check.py --restart
