# deploy/watchdog.timer
# live-scheduler-deploy FR-09 — 워치독 타이머 (15분마다 검사)
[Unit]
Description=auto-stock watchdog timer

[Timer]
OnBootSec=5min
OnActiveSec=15min
OnUnitActiveSec=15min
Unit=watchdog.service

[Install]
WantedBy=timers.target
